# rootfs/etc/s6-overlay/s6-rc.d/fr24uat-feed/run
#!/bin/sh
exec /etc/s6-overlay/scripts/fr24uat-feed
