# etc/services.d/fr24feed_log/run
#!/usr/bin/with-contenv sh
#shellcheck shell=sh

exec \
  tail -F /var/log/fr24feed.log
