# etc/services.d/fr24feed_log/run
#!/usr/bin/with-contenv sh

tail -F /var/log/fr24feed.log
