# rootfs/etc/services.d/fr24feed/run
#!/usr/bin/with-contenv sh
#shellcheck shell=sh

exec \
  /usr/bin/fr24feed > /dev/null 2>&1
