# etc/services.d/fr24feed/run
#!/usr/bin/with-contenv sh

/usr/bin/fr24feed 2>&1 > /dev/null
